# supervisord.conf
[program:twistd-web]
command = twistd -n --pidfile= web --port=tcp:8080 --class=hello.resource
process_name = twistd-web
directory = /home/dotcloud/current/
priority = 10
redirect_stderr = true
